# armonic-agent-socket.service
[Unit]
Description=Armonic Socket Agent

[Service]
ExecStart=/usr/bin/armonic-agent-socket -vvv
Restart=on-abort

[Install]
WantedBy=default.target
